# crates/config/tests/fixtures/service/launchd_spaces.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.toolu.runner.a &amp; b.repo</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/dev name/apps/toolu &amp; runner/toolu-runner</string>
    <string>run</string>
    <string>--config</string>
    <string>/Users/dev name/.toolu-runner/runners/a &amp; b/repo/config.toml</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/dev name/.toolu-runner/runners/a &amp; b/repo/_diag/service.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/dev name/.toolu-runner/runners/a &amp; b/repo/_diag/service.err.log</string>
</dict>
</plist>
